# certbot/certbot-docker.service
[Unit]
Description=certbot 定时任务


[Service]
ExecStart=/usr/bin/docker run -i --rm -v /etc/letsencrypt/:/etc/letsencrypt/ certbot:latest renew

